# pairdrop/rootfs/etc/s6-overlay/s6-rc.d/svc-pairdrop/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
exec sleep 5
