# pairdrop/rootfs/etc/s6-overlay/s6-rc.d/svc-pairdrop/finish
#!/usr/bin/with-contenv bash
echo "[pairdrop] === Finish: waiting for port 3000 to be released ==="
for i in $(seq 1 30); do
    if ! nc -z localhost 3000 2>/dev/null; then
        echo "[pairdrop] Port 3000 freed after ${i}s"
        break
    fi
    echo "[pairdrop] Port 3000 still in use... (${i}/30)"
    sleep 1
done
echo "[pairdrop] === Finish done ==="
