# pairdrop/rootfs/etc/s6-overlay/s6-rc.d/svc-pairdrop/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set +u

echo "=== PAIRDROP DEBUG (PID=$$, PPID=$(ps -o ppid= -p $$ 2>/dev/null || echo unknown)) ==="

echo "=== ALL PROCESSES ==="
ps aux 2>/dev/null || ps -ef 2>/dev/null || echo "ps failed"

echo "=== LISTENING PORTS ==="
ss -tlnp 2>/dev/null || netstat -tlnp 2>/dev/null || echo "ss/netstat not available"

echo "=== PORT 3000 ==="
fuser 3000/tcp 2>/dev/null && echo "Port 3000 in use by above PIDs" || echo "Port 3000 free"

echo "=== S6 SERVICES ==="
ls -la /etc/s6-overlay/s6-rc.d/ 2>/dev/null || true
for d in /etc/s6-overlay/s6-rc.d/*/; do
    echo "--- $d ---"
    cat "${d}type" 2>/dev/null || echo "(no type)"
    ls -la "${d}" 2>/dev/null || true
done

echo "=== /etc/services.d ==="
ls -laR /etc/services.d/ 2>/dev/null || echo "no services.d"

echo "=== RUN FILES ==="
for f in /etc/s6-overlay/s6-rc.d/*/run; do
    echo "--- $f ---"
    head -3 "$f" 2>/dev/null || true
done

echo "=== ha_entrypoint.sh shebang detection ==="
cat /etc/s6-overlay/s6-rc.d/svc-pairdrop/run 2>/dev/null || true

echo "=== SLEEPING 120s (no node start) ==="
sleep 120
echo "=== DONE ==="
