# pairdrop/rootfs/etc/s6-overlay/s6-rc.d/svc-pairdrop/run
#!/usr/bin/with-contenv bashio
set +u
source /usr/local/lib/bashio-standalone.sh

echo "=== PAIRDROP DIAGNOSTIC MODE ==="
echo "Hostname: $(hostname)"
echo "Date: $(date)"
echo "User: $(id)"
echo ""
echo "=== Processes ==="
ps aux 2>/dev/null || ps -ef 2>/dev/null
echo ""
echo "=== Port 3000 status ==="
ss -tlnp 2>/dev/null | grep 3000 || netstat -tlnp 2>/dev/null | grep 3000 || echo "(no port tool)"
nc -z localhost 3000 2>/dev/null && echo "PORT 3000 IN USE" || echo "PORT 3000 FREE"
echo ""
echo "=== Network interfaces ==="
ip addr 2>/dev/null || ifconfig 2>/dev/null
echo ""
echo "=== Relevant env vars ==="
env | grep -iE "PORT|NODE|NPM|RATE|WS_F|HOST|S6|HOME|PATH|DEBUG|PAIR" || true
echo ""
echo "=== /app/pairdrop exists? ==="
ls -la /app/pairdrop/ 2>/dev/null || echo "NO /app/pairdrop"
echo ""
echo "=== /config contents ==="
ls -la /config/ 2>/dev/null || echo "NO /config"
echo ""
echo "=== Run script path ==="
readlink -f "$0" 2>/dev/null || echo "(readlink not available)"
echo ""
echo "=== Sleeping forever (diagnostic mode) ==="
sleep infinity
