# pairdrop/rootfs/etc/s6-overlay/s6-rc.d/svc-pairdrop/run
#!/usr/bin/with-contenv bash
echo "================================================"
echo "CUSTOM RUN SCRIPT v1.11.2-30-marker"
echo "================================================"

LOCATION="${LOCATION:-/config}"
export PORT=8080

cd /app/pairdrop

OPT_RATE_LIMIT=""
OPT_WS_FALLBACK=""
if [[ "${RATE_LIMIT:-}" = "true" ]]; then OPT_RATE_LIMIT="--rate-limit"; fi
if [[ "${WS_FALLBACK:-}" = "true" ]]; then OPT_WS_FALLBACK="--include-ws-fallback"; fi

echo "Starting PairDrop on PORT=${PORT} LOCATION=${LOCATION}"

if [[ -z "${LSIO_NON_ROOT_USER:-}" ]]; then
    exec s6-setuidgid abc node server/index.js ${OPT_RATE_LIMIT} ${OPT_WS_FALLBACK}
else
    exec node server/index.js ${OPT_RATE_LIMIT} ${OPT_WS_FALLBACK}
fi
