# rezuscloud/rootfs/etc/services.d/rezuscloud/finish
#!/usr/bin/execlineb -S
s6-svscanctl -t /run/service
