# pairdrop/rootfs/etc/services.d/pairdrop/run
#!/usr/bin/env bashio
set +u

echo "================================================"
echo "CUSTOM RUN SCRIPT v1.11.2-29 PORT=8080"
echo "================================================"

bashio::log.info "Starting PairDrop add-on"

LOCATION="$(bashio::config 'data_location')"
if [[ "$LOCATION" = "null" || -z "$LOCATION" ]]; then LOCATION=/config; fi
DEBUG_MODE="$(bashio::config 'debug_mode')"
export DEBUG_MODE
export PORT=8080

mkdir -p "$LOCATION"
cd /app/pairdrop

RATE_LIMIT="$(bashio::config 'rate_limit')"
WS_FALLBACK="$(bashio::config 'ws_fallback')"

OPT_ARGS=""
if [[ "${RATE_LIMIT:-}" = "true" ]]; then OPT_ARGS="$OPT_ARGS --rate-limit"; fi
if [[ "${WS_FALLBACK:-}" = "true" ]]; then OPT_ARGS="$OPT_ARGS --include-ws-fallback"; fi

bashio::log.info "Data location: ${LOCATION}"
bashio::log.info "PORT=${PORT}"
bashio::log.info "Starting PairDrop server..."

exec node server/index.js $OPT_ARGS
